# packaging/fangd.service
[Unit]
Description=Fang daemon — performance and fan control for Razer Blade laptops
Documentation=https://github.com/bladeandsoulx/fang-razer-linux

[Service]
Type=exec
ExecStart=/usr/bin/fangd
ExecStopPost=/usr/bin/fangd --restore-auto
Restart=on-failure
RestartSec=3
TimeoutStopSec=8
StateDirectory=fangd
# Hardening: the daemon only needs hidraw, hwmon/nvml and its state dir.
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
